# analytics-service.service
[Unit]
Description=WMP Analytics Service
After=network.target

[Service]
Type=simple
User=appuser
WorkingDirectory=/app
ExecStart=/usr/local/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 4
Restart=on-failure
RestartSec=10

Environment=DATABASE_URL=postgresql+asyncpg://analytics_svc_user:localdev123@postgresql-dev.raghudevopsb88.online:5432/wmp
Environment=DB_SCHEMA=analytics_schema
Environment=LOG_LEVEL=INFO
Environment=ENVIRONMENT=production

[Install]
WantedBy=multi-user.target